# zabbix-agent2/rootfs/etc/services.d/zabbix-agent2/run
#!/usr/bin/with-contenv bash
# Dit script draait als één van de s6-services
exec /usr/bin/run.sh
