# zabbix-agent2/rootfs/etc/services.d/zabbix-agent2/run
#!/usr/bin/with-contenv bashio
# s6-service wrapper; draait nรก de s6-init
exec /usr/local/bin/run.sh
